# production/cron/crontabs/crontab-deletion
5 * * * *    /app/scripts/deactivate-projects.sh > /proc/1/fd/1 2>&1
15 * * * *   /app/scripts/delete-users.sh > /proc/1/fd/1 2>&1
20 * * * *   /app/scripts/delete-projects.sh > /proc/1/fd/1 2>&1
